# container-files/etc/supervisor.d/crond.conf
[program:crond]
command = /usr/sbin/crond -n -m off
autorestart = true