# container-files/etc/supervisor.d/filebeat.conf
[program:filebeat]
command = /usr/bin/filebeat -e -c "/etc/filebeat/filebeat.yml"
autorestart = true