# container-files/etc/supervisor.d/haproxy.conf
[program:haproxy]
command = /usr/local/sbin/haproxy -f /usr/local/etc/haproxy/haproxy.cfg %(ENV_DEBUG)s
autorestart = true