# container-files/etc/supervisor.d/modsecurity.conf
[program:modsecurity]
command = /usr/local/bin/modsecurity -f /usr/local/etc/modsecurity/modsecurity.conf -n %(ENV_MODSEC_NUM_WORKERS) %(ENV_DEBUG)s
autorestart = true